# .config/containers/systemd/esphome.container
# Podman Quadlet container for esphome
# Generated by scaffold-homelab-batch.py

[Unit]
Description=esphome container
After=network-online.target
Wants=network-online.target

[Container]
ContainerName=esphome
Image=docker.io/esphome/esphome:latest

Network=host

EnvironmentFile=homelab/esphome/esphome.env
Environment=TZ=Europe/Brussels

WorkingDirectory=~
Volume=homelab/esphome/config:/config:z

AutoUpdate=registry

Restart=always

[Install]
WantedBy=default.target
